# deploy_tools/supervisord.template.conf
[program:SITENAME]
command=/etc/gunicorn/SITENAME_start
redirect_stderr=true
stdout_logfile=/var/log/gunicorn/SITENAME.log
enviroment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8
